# crontab
* * * * * node /app/src/cron_runner.js >> /cron/cron.log 2>&1